# supervisord/supervisord-openvpn.conf
[program:openvpn]
command=/usr/local/bin/run-openvpn.sh
autorestart=false
startretries=0
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
