# supervisord/supervisord-db1000n.conf
[program:db1000n]
command=/usr/local/bin/db1000n.sh
autorestart=false
startretries=0
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
